# package/base-files/files/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

#系统内存优化

if [ -f /sbin/memory_release ];then
        /sbin/memory_release &
fi


if [ -f /sbin/boot_times.sh ];then
        /sbin/boot_times.sh &
fi

exit 0
